# src/supervisord.conf
[inet_http_server]
port=127.0.0.1:9001


[program:csnet_client_darwin_amd64]
command = /Users/jcleng/Downloads/csnet_client_darwin_amd64
